# server-config/backend-api.service
[Unit]
Description=AgentCloud Backend API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/agentcloud
Environment="NODE_ENV=production"
EnvironmentFile=-/opt/agentcloud/.env
# RELAY_INTERNAL_SECRET should be set in /opt/agentcloud/.env to secure internal endpoints
ExecStart=/usr/bin/npx tsx backend/src/server.ts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
