# server-config/backend-api.service
[Unit]
Description=AgentCloud Backend API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/agentcloud
Environment="NODE_ENV=production"
EnvironmentFile=/opt/agentcloud/.env
ExecStart=/usr/bin/node --loader tsx/esm backend/src/server.ts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

